# templates/etc/systemd/system/data-container.service.erb
[Unit]
Description=Docker Data Container for <%= @title %>
Requires=docker.service <%= scope.function_systemd_requires([{:depends=>@depends, :systemd_depends=>@systemd_depends}]) %>
After=docker.service <%= scope.function_systemd_requires([{:depends=>@depends, :systemd_depends=>@systemd_depends}]) %>

[Service]
Type=oneshot
Restart=no
RemainAfterExit=yes
<% if @systemd_env_file %><%= "EnvironmentFile=#{@systemd_env_file}" %><% end %>
<% if @pull_image %>ExecStartPre=/usr/bin/docker pull <%= @image %><% end %>
ExecStartPre=-/usr/bin/docker run \
    <%= @docker_run_options %> \
    <%= @image %>
ExecStart=/usr/bin/docker start <%= @title %>
ExecStop=/usr/bin/docker stop <%= @title %>

[Install]
WantedBy=multi-user.target
