# templates/etc/systemd/system/run-container.service.erb
[Unit]
Description=Docker Container Service for <%= @title %>
Requires=docker.service <%= scope.function_systemd_requires([{:depends=>@depends, :systemd_depends=>@systemd_depends}]) %>
After=docker.service <%= scope.function_systemd_requires([{:depends=>@depends, :systemd_depends=>@systemd_depends}]) %>

[Service]
Type=simple
Restart=always
RestartSec=5
<% if @systemd_env_file %><%= "EnvironmentFile=#{@systemd_env_file}" %><% end %>
ExecStartPre=-/usr/bin/docker stop <%= @title %>
ExecStartPre=-/usr/bin/docker rm <%= @title %>
<% if @pull_image %>ExecStartPre=/usr/bin/docker pull <%= @image_arg %><% end %>
ExecStart=/usr/bin/docker run --rm \
    <%= @docker_run_options %> \
    <%= @image_arg %> <%= @command %>
ExecStop=/usr/bin/docker stop <%= @title %>

[Install]
WantedBy=multi-user.target
